# supervisor/conf.d/bot.conf
[program:ircbot]
command=/usr/bin/npm run-script deploy
directory=/opt/bot/
autostart=true
autorestart=true